# deploy/dogbook-backup.timer
[Unit]
Description=Daily Dogbook Database Backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
