# deploy/dogbook-backup.service
[Unit]
Description=Backup Dogbook Database

[Service]
Type=oneshot
User=caddy
ExecStart=/bin/bash -c '/usr/bin/sqlite3 /srv/dogbook/data/keystone.db ".backup /srv/dogbook/backups/keystone-$(date +%%Y%%m%%d).db"'
ExecStartPost=/usr/bin/find /srv/dogbook/backups -name "keystone-*.db" -mtime +30 -delete
