# deploy/dogbook.service
[Unit]
Description=Dogbook Backend (KeystoneJS)
After=network.target

[Service]
Type=simple
User=dogbook
WorkingDirectory=/srv/dogbook/current/backend
ExecStart=/home/dogbook/.nvm/versions/node/v22.21.1/bin/node /srv/dogbook/current/backend/node_modules/.bin/keystone start
Environment=NODE_ENV=production
Environment=PATH=/home/dogbook/.nvm/versions/node/v22.21.1/bin:/usr/local/bin:/usr/bin:/bin
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
